# metropi.service
[Unit]
Description=Metro Pi Display Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/metropi
ExecStart=/home/pi/metropi/venv/bin/python /home/pi/metropi/main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
